# openlitespeed/openlitespeed.conf
[program:openlitespeed]
command=/usr/local/lsws/bin/openlitespeed -d
autostart=true
autorestart=true